# iso/airootfs/etc/systemd/system/greetd.service
[Unit]
Description=Greetd login manager
Documentation=man:greetd(1)
After=systemd-user-sessions.service plymouth-quit-wait.service
After=getty@tty1.service

[Service]
Type=idle
ExecStart=/usr/bin/greetd
Restart=always

[Install]
Alias=display-manager.service
